# message-broker/carbonio-message-broker.service
[Unit]
Description=Carbonio message broker
Requires=network-online.target
After=network.target epmd@0.0.0.0.socket
Wants=network.target epmd@0.0.0.0.socket

[Service]
Type=notify
User=carbonio-message-broker
Group=carbonio-message-broker
UMask=0027
NotifyAccess=all
TimeoutStartSec=600
LimitNOFILE=32768
Restart=on-failure
RestartSec=10
WorkingDirectory=/var/lib/rabbitmq
ExecStart=/opt/zextras/common/lib/rabbitmq/bin/rabbitmq-server
ExecStop=/opt/zextras/common/lib/rabbitmq/bin/rabbitmqctl shutdown
SuccessExitStatus=69

[Install]
WantedBy=multi-user.target
